# tools/hakpak-fulfillment.service
[Unit]
Description=HakPak Pro License Fulfillment Service
After=network.target

[Service]
Type=notify
User=hakpak-vendor
Group=hakpak-vendor
WorkingDirectory=/secure/hakpak-vendor/tools
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/secure/hakpak-vendor/.env
ExecStart=/usr/local/bin/gunicorn --workers 4 --bind 0.0.0.0:5000 --timeout 60 --access-logfile /var/log/hakpak_access.log --error-logfile /var/log/hakpak_error.log stripe_webhook:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/secure/hakpak-vendor/licenses /var/log
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
